# systemd/acidzero.service
[Unit]
Description=Acid Zero pentest handheld UI (TFT)
After=pwnagotchi.service
[Service]
# Detach the kernel framebuffer console (fbcon) + disable blanking so the Linux
# console / getty / cursor never draw over the UI (the app owns /dev/fb0 directly).
ExecStartPre=/usr/local/bin/acid-console-off.sh
ExecStart=/home/pi/.pwn/bin/python3 /usr/local/bin/acidzero.py
Restart=always
RestartSec=10
[Install]
WantedBy=multi-user.target
